# installer/templates/blackbox-models.service
# /etc/systemd/system/blackbox-models.service
# Generated by BlackBox blackbox-install-localstack.sh — DO NOT EDIT BY HAND.
#
# llama-swap front door on 127.0.0.1:9098 supervising the on-box local-model
# stack: the "retrieval" group (embed-qwen3-8b + rerank-qwen3-8b) and the
# "audio" group (speaches STT + our qwen-tts server). One binary spawns/kills
# the member servers on demand per the generated
# LOCALSTACK_HOME_PLACEHOLDER/llama-swap-config.yaml. --watch-config
# auto-restarts the whole proxy on any config edit (§6).
#
# REAL_USER_PLACEHOLDER / REAL_HOME_PLACEHOLDER / LOCALSTACK_BIN_PLACEHOLDER /
# LOCALSTACK_HOME_PLACEHOLDER are substituted at install time (same
# sed-template flow as vllm-reranker.service / zellij-web.service).
#
# GPU and CPU boxes both run this unit; only the member builds/models differ
# (§7). Members lazy-load on first request, so the front door answers /health
# with ZERO weights resident — weights download later in the wizard
# (disk-gated, POST /local-models/download). LD_LIBRARY_PATH points at the
# localstack bin dir so a CUDA llama-server prebuilt finds its bundled libs.

[Unit]
Description=BlackBox local model stack (llama-swap front door :9098)
After=network-online.target
Wants=network-online.target

[Service]
User=REAL_USER_PLACEHOLDER
WorkingDirectory=REAL_HOME_PLACEHOLDER
Environment=LD_LIBRARY_PATH=LOCALSTACK_BIN_PLACEHOLDER
ExecStart=LOCALSTACK_BIN_PLACEHOLDER/llama-swap --config LOCALSTACK_HOME_PLACEHOLDER/llama-swap-config.yaml --listen 127.0.0.1:9098 --watch-config
Restart=on-failure
RestartSec=15
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
