# installer/templates/vllm-reranker.service
# /etc/systemd/system/vllm-reranker.service
# Generated by BlackBox blackbox-install-reranker.sh — DO NOT EDIT BY HAND.
#
# vLLM serving Qwen/Qwen3-Reranker-0.6B on port 8091 for the M11 cross-encoder
# rerank stage (Orchestrator/rerank.py — the activation checklist lives in its
# module docstring). Templated verbatim from the deployed MS-02 Ultra unit
# (the proven configuration); the start script owns the vllm incantation
# (port, --gpu-memory-utilization 0.20 so the Ollama embedder can share the
# card, the Qwen3 sequence-classification hf-overrides).
#
# GPU boxes ONLY: blackbox-install-reranker.sh gates on nvidia-smi presence
# + >=8000 MB VRAM and never installs this unit on a CPU box.
#
# REAL_USER_PLACEHOLDER / REAL_HOME_PLACEHOLDER are substituted at install
# time with the actual system user/home (same user that owns blackbox.service
# — mirrors zellij-web.service's parameterization).
#
# TimeoutStartSec=600: the first start downloads the model weights from
# Hugging Face (~1.2 GB) before the engine can bind the port.

[Unit]
Description=vLLM Qwen3-Reranker 0.6B (BlackBox M11 reranker)
After=network-online.target
Wants=network-online.target

[Service]
User=REAL_USER_PLACEHOLDER
WorkingDirectory=REAL_HOME_PLACEHOLDER
ExecStart=REAL_HOME_PLACEHOLDER/start-reranker.sh
Restart=on-failure
RestartSec=15
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
