# installer/templates/zellij-web.service
# /etc/systemd/system/zellij-web.service
# Generated by BlackBox install.sh — DO NOT EDIT BY HAND.
#
# Runs the Zellij web client daemon that backs the AI BlackBox CLI Agent
# terminal bridge. The orchestrator (blackbox.service) mints per-session
# tokens against this daemon's token store (~/.local/share/zellij/tokens.db)
# and the Portal embeds the daemon's web UI per session.
#
# REAL_USER_PLACEHOLDER is substituted at install time with the actual
# system user (same user that owns blackbox.service — Phase 0 audit C0).
# Shared user namespace lets the orchestrator read tokens.db with no sudo.
#
# Type=simple (NOT forking with --daemonize): Phase 0 finding #3 confirmed
# `zellij web --daemonize` is broken in 0.44.3 — the fork claims success but
# the child silently dies. systemd's own Type=simple foregrounding works.
#
# Port: configured via ~/.config/zellij/config.kdl (written by install.sh),
# NOT via a CLI flag. Default = 9097 (see installer/templates/zellij-version).
#
# NO PrivateTmp / ProtectSystem / ProtectHome hardening: per memory
# `cli_agent_systemd_hardening.md`, these silently break Zellij's socket
# (/tmp) + session state (~/.local/share/zellij/) requirements. The unit
# inherits the shared $REAL_USER trust boundary already established by
# blackbox.service; adding sandboxing here without re-architecting the
# token-store handoff would break the CLI Agent feature on first boot.

[Unit]
Description=AI BlackBox CLI Agent — Zellij web client daemon
Documentation=https://github.com/zellij-org/zellij
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REAL_USER_PLACEHOLDER
ExecStart=/usr/local/bin/zellij web
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zellij-web

[Install]
WantedBy=multi-user.target
