# MCP/deploy/blackbox-mcp.service
# BlackBox MCP HTTP Server -- persistent Streamable-HTTP (remote) transport.
# Install on a fresh box:
#   sudo cp MCP/deploy/blackbox-mcp.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now blackbox-mcp.service
# The absolute paths below MUST match the repo's install location.
[Unit]
Description=BlackBox MCP HTTP Server (remote tool server, bearer/OAuth authed)
After=network.target blackbox.service
Wants=blackbox.service

[Service]
Type=simple
User=ai-black-box-fc
WorkingDirectory=/home/ai-black-box-fc/Desktop/blackbox_poc./blackbox_poc
Environment=BLACKBOX_MCP_TRANSPORT=http
Environment=BLACKBOX_MCP_HTTP_HOST=127.0.0.1
Environment=BLACKBOX_MCP_HTTP_PORT=9093
Environment=BLACKBOX_ROOT=/home/ai-black-box-fc/Desktop/blackbox_poc./blackbox_poc
Environment=BLACKBOX_URL=http://localhost:9091
Environment=BLACKBOX_MCP_TOKENS_FILE=/home/ai-black-box-fc/Desktop/blackbox_poc./blackbox_poc/Manifest/mcp_tokens.json
# For OAuth clients (claude.ai / native Claude Desktop), set THIS box's public
# Funnel origin so the OAuth discovery metadata points at your host. There is no
# default (a hard-coded host would mis-point every other install); OAuth discovery
# returns 503 until set, and bearer-token auth is unaffected:
# Environment=BLACKBOX_MCP_PUBLIC_URL=https://<your-host>.ts.net:8443
# Roster enforcement (reject any token bound to a non-roster operator at startup)
# is ON by default; set BLACKBOX_MCP_ROSTER_ENFORCE=0 only to disable it.
ExecStart=/home/ai-black-box-fc/Desktop/blackbox_poc./blackbox_poc/MCP/venv/bin/python /home/ai-black-box-fc/Desktop/blackbox_poc./blackbox_poc/MCP/blackbox_mcp_server.py --transport http
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
